# s6-rc.d/rsyslog/run
#!/bin/sh
. /etc/s6-overlay/s6-rc.d/.config-init

mkdir -p /var/spool/rsyslog

exec /usr/sbin/rsyslogd -n
