# s6-rc.d/wazuh-agent/finish
#!/command/execlineb -P
/var/ossec/bin/wazuh-control stop