# s6-rc.d/wazuh-agent/run
#!/bin/sh
. /etc/s6-overlay/s6-rc.d/.config-init
/var/ossec/bin/wazuh-control start
exec sleep infinity
